# src/services/mysql/run
docker run -d -P --volumes-from=mariadb-data --name=mariadb market/mysql
